# conf/systemd.service.in
# This file is part of package @PACKAGE_NAME@ @PACKAGE_VERSION@
#
# Reference:
#
#	http://www.freedesktop.org/software/systemd/man/systemd.service.html
#
#
[Unit]
Description="PW3270 IPC Service Daemon"
Requires=var-run.mount nss-lookup.target network.target
After=var-run.mount nss-lookup.target network.target

[Service]
Type=notify
ExecStartPre=/bin/echo 'Starting @PACKAGE_NAME@ IPC Service'
ExecStart=/usr/sbin/@PACKAGE_NAME@ --core --pidfile=/var/run/@PACKAGE_NAME@.pid
Restart=on-failure
PIDFile=/var/run/@PACKAGE_NAME@.pid
RestartSec=10s
WatchdogSec=120

[Install]
WantedBy=multi-user.target

